# docs/ML/hermes.service.example
[Unit]
Description=Hermes Autonomous Infrastructure Agent
After=docker.service network-online.target
Requires=docker.service

[Service]
Restart=always
RestartSec=10
ExecStartPre=-/usr/bin/docker stop hermes
ExecStartPre=-/usr/bin/docker rm hermes
ExecStart=/usr/bin/docker run --name hermes \
  --network=host \
  --security-opt apparmor=unconfined \
  --security-opt seccomp=unconfined \
  --env-file /opt/hermes/.env \
  -v /opt/hermes/data:/opt/data \
  -v /usr/sbin:/usr/sbin:ro \
  -v /usr/bin:/usr/bin:ro \
  -v /usr/local/bin/kubectl:/usr/local/bin/kubectl:ro \
  -v /home/ksh/infra:/infra \
  -v /home/ksh/.ssh:/home/ksh/.ssh:ro \
  -e "HOME=/home/ksh" \
  nousresearch/hermes-agent:latest \
  gateway run
ExecStop=/usr/bin/docker stop hermes

[Install]
WantedBy=multi-user.target
